# hosts/shootingstar/home/.config/systemd/user/ngrok.service
[Unit]
Description=ngrok ssh

[Service]
WorkingDirectory=/home/nox/
ExecStart=/home/nox/.local/bin/ngrok tcp 22
Restart=always

[Install]
WantedBy=default.target
