# home/.config/systemd/user/gmail-notifications-work.timer
[Unit]
Description=Timer for gmail unread notifications.
Requires=gmail-notifications-work.service

[Timer]
OnCalendar=*-*-* *:*:00
AccuracySec=1s

[Install]
WantedBy=timers.target
